# templates/debian-bullseye-mate/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -a dodo --noclear %I $TERM
